# visionect-v3-allinone/rootfs/etc/services.d/postgres/run
#!/command/with-contenv bash
exec /usr/lib/postgresql/*/bin/postgres -D "${PGDATA}"
