# visionect-v3-allinone/rootfs/etc/services.d/redis/run
#!/usr/bin/execlineb -P
fdmove -c 2 1
s6-setuidgid root
mkdir -p /data/redis
exec redis-server --port 6379 --dir /data/redis --save "" --appendonly no